# galaxy-book-audio/ionsoundsleep.service
[Unit]
Description=Custom Sound Sleep Fix

[Service]
Type=idle
ExecStart=/opt/TO912.sh

[Install]
WantedBy=suspend.target
