# galaxy-book-audio/ionsound.service
[Unit]
Description=Custom Sound

[Service]
Type=idle
ExecStart=/opt/TO912.sh

[Install]
WantedBy=multi-user.target
